# autobuilder.service.in
[Unit]
Description=distri autobuilder

[Service]
ExecStart=@AUTOBUILDER@

[Install]
WantedBy=default.target
