# distri-checkupstream.timer
[Unit]
Description=Run distri-checkupstream daily

[Timer]
OnCalendar=*-*-* 07:50:00

[Install]
WantedBy=timers.target
